# lis-private-backup/var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.qO7scq/crontab installed on Thu May 21 00:29:14 2015)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
00 * * * * /root/correrbackup.sh
